# assets/prometheus-modem-exporter.service
[Unit]
Description=Prometheus exporter for modem metrics
Requires=network-online.target
After=network-online.target

[Service]
DynamicUser=yes
User=consul
Restart=on-failure
ExecStart=/usr/local/bin/modem-exporter
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
